# sysutils/sensu-go/files/sensu-agent.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: sensu-agent
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable sensu-agent:
#
# sensu_agent_enable (bool):	Set to YES to enable sensu-agent
# 				Default: NO
# sensu_agent_flags (str):	extra flags to agent
#				Default: ""
# sensu_agent_env (str):	list of environment variables to pass
#				Default: HOME=/home/sensu

. /etc/rc.subr

name=sensu_agent
rcvar=sensu_agent_enable

load_rc_config ${name}

: ${sensu_agent_enable:="NO"}
: ${sensu_agent_restart_delay:="10"}
: ${sensu_agent_env:="HOME=/home/sensu"}

pidfile="/var/run/sensu-agent.pid"
command=/usr/sbin/daemon
sensu_agent_command=%%PREFIX%%/bin/sensu-agent
command_args="-rP ${pidfile} -S -R ${sensu_agent_restart_delay} \
    -T ${name} ${sensu_agent_command} \
    start ${sensu_agent_flags}"

load_rc_config $name
run_rc_command "$1"
