# sysutils/sensu-go/files/sensu-backend.in
#!/bin/sh
# $FreeBSD$

# PROVIDE: sensu-backend
# REQUIRE: LOGIN
# KEYWORD: shutdown

# Add the following lines to /etc/rc.conf to enable sensu-backend:
#
# sensu_backend_enable (bool):	Set to YES to enable sensu-backend
# 				Default: NO
# sensu_backend_flags (str):	Extra flags to backend
# 				Default: ""
# sensu_backend_env (str):	list of environment variables to pass
#				Default: HOME=/home/sensu

. /etc/rc.subr

name=sensu_backend
rcvar=sensu_backend_enable

load_rc_config ${name}

: ${sensu_backend_enable:="NO"}
: ${sensu_backend_flags:=""}
: ${sensu_backend_restart_delay:="10"}
: ${sensu_backend_env:="HOME=/home/sensu"}

pidfile="/var/run/sensu-backend.pid"
command=/usr/sbin/daemon
sensu_backend_command=%%PREFIX%%/bin/sensu-backend
command_args="-rP ${pidfile} -S -R ${sensu_backend_restart_delay} \
    -T ${name} ${sensu_backend_command} \
    start ${sensu_backend_flags}"

load_rc_config $name
run_rc_command "$1"
